# deployment/ansible/roles/pytrader_app/templates/pytrader-scrape.timer.j2
[Unit]
Description=Run PyTrader Data Scraper daily

[Timer]
OnCalendar={{ pytrader_scrape_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
